# deploy/yijunbot.service
[Unit]
Description=Yijun Discord bot
After=network-online.target llama-server.service
Wants=network-online.target
Requires=llama-server.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/Yijun.skill
EnvironmentFile=/etc/yijunbot/.env
ExecStart=/home/ubuntu/yijunbot-venv/bin/python -m bot.discord_bot
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
